# imageroot/systemd/user/hermes-auth@.service
[Unit]
Description=Hermes dashboard auth proxy %i
Requires=hermes-pod@%i.service hermes@%i.service hermes-dashboard@%i.service
After=hermes-pod@%i.service hermes@%i.service hermes-dashboard@%i.service
PartOf=hermes@%i.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-%S/state/environment
EnvironmentFile=-%S/state/agent_%i.env
WorkingDirectory=%S/state
Restart=on-failure
KillMode=none
Type=notify
NotifyAccess=all
ExecStartPre=/usr/local/bin/runagent sync-agent-runtime --agent-id %i
ExecStartPre=-/usr/bin/podman rm --force hermes-auth-%i
ExecStart=runagent /usr/bin/podman run \
		--name hermes-auth-%i \
		--replace \
		--rm \
		--sdnotify=conmon \
		--cgroups=no-conmon \
		--pod hermes-pod-%i \
		--env-file %S/state/agent_%i.env \
		--env-file %S/state/agent_%i_secrets.env \
		--env BASE_URL=/hermes-%i \
		--env DASHBOARD_UPSTREAM_URL=http://127.0.0.1:9120 \
		--env AUTH_PROXY_PORT=9119 \
		--tz=${TIMEZONE} \
		${HERMES_AGENT_AUTH_IMAGE} \
		python /app/authproxy.py
ExecStop=/usr/bin/podman stop --ignore --time 10 hermes-auth-%i
ExecStopPost=-/usr/bin/podman rm --force hermes-auth-%i
TimeoutStopSec=70